# update-motd.d/10-sysinfo
#!/bin/sh

DATE="Today is:\t\t$(date)"
UNAME="Kernel information:\t$(uname -srm)"

echo $DATE
echo $UNAME
